# com.useful_appcoll.docket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.useful_appcoll.docket</string>

  <!-- Update path to your script -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/your_username/path/to/useful_appcoll/run_docket.sh</string>
  </array>

  <!-- 9:00 AM local time, Monday–Friday (launchd weekdays: 1=Mon, ..., 5=Fri) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <!-- Update path to your log folder -->
  <key>StandardOutPath</key>
  <string>/Users/your_username/path/to/useful_appcoll/logs/launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/your_username/path/to/useful_appcoll/logs/launchd_stderr.log</string>
</dict>
</plist>
